# 2018-15-sesja-linuksowa-oauth-nginx/demo-4/oauthsub.service
[Unit]
Description=oauthsub service
After=nginx.service

[Service]
Type=simple
ExecStart=/usr/bin/oauthsub -c /etc/oauthsub.py
User=nginx
Restart=on-abort

[Install]
WantedBy=multi-user.target