# 2016-anonimowi-admini-systemd-workshop/sd-notify-python/systemd-notify-app.socket
[Unit]
Description=Socket activation for simple systemd-notify app

[Socket]
ListenStream=1025

[Install]
WantedBy=sockets.target
