# 2017-14-sesja-linuksowa-systemd/journal-python/journal-app-service.service
[Unit]
Description=Simple Journal writer application

[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/python3 /home/docent/git/conferences/systemd-devops-meetup-wroclaw-2016/journal-python/journal-app-service.py

[Install]
WantedBy=multi-user.target